# download_organizer.service
[Unit]
Description=download organizer will sort files into respective directories

[Service]
Type=oneshot
ExecStart=/usr/bin/download_organizer
